# build/Lede_opewrtR3G/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

cat /dev/null > /etc/bench.log
echo " (CpuMark : 187821.997633" >> /etc/bench.log
echo " Scores)" >> /etc/bench.log

if [ -f "/etc/bench.log" ]; then
 sed -i '/coremark/d' /etc/crontabs/root
 crontab /etc/crontabs/root
fi

exit 0
